# recipes-bsp/optee/files/optee.service
[Unit]
Description=OP-TEE Supplicant Service

[Service]
ExecStart=/usr/bin/tee-supplicant

[Install]
WantedBy=multi-user.target

