# recipes-core/tc-enable-removable-disk/tc-enable-removable-disk/enable-removable-disk.path
[Unit]
Description=This service to enable removable disk after application ready.

[Path]
PathExists=/var/run/MediaPlayer.pid

[Install]
WantedBy=graphical.target
